# launchd/ai.xiaodao.ime.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.xiaodao.ime</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/chengang/xiaodao-ime/.venv/bin/python</string>
        <string>/Users/chengang/xiaodao-ime/app.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/chengang/xiaodao-ime</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/chengang/xiaodao-ime/logs/launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/chengang/xiaodao-ime/logs/launchd.err.log</string>

    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
